# application/config/phorge-sshd.sv.conf
[program:phorge-sshd]
command=/usr/sbin/sshd -D -f /etc/phorge-ssh/sshd_config.phorge