# application/config/php7-fpm.sv.conf
[program:php7.4-fpm]
command=php-fpm7.4 -F