# .devcontainer/application/config/aphlict.sv.conf
[program:aphlict]
command=/srv/phorge/phorge/bin/aphlict debug --config /srv/phorge/phorge/conf/aphlict/aphlict.default.json
user=phorge-daemon